# nim/supervisord.conf
[program:fserve]
command=sh -c 'sleep 10 && /src/main'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[supervisord]
